# cron.d/chronograph
* * * * * www-data bash /sites/oserver/deploy/chronograph.sh /sites/oserver
